# templates/default/etc/systemd/system/pgbouncer.service.erb
[Unit]
Description=pgbouncer-<%= @db_alias %>
After=network.target

[Service]
User=<%= @user %>
ExecStart=/usr/sbin/pgbouncer /etc/pgbouncer/pgbouncer-<%= @db_alias %>.ini
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
